# etc/systemd/system/webhooks.service
[Unit]
Description=Webhooks

[Service]
ExecStart=/bin/bash /apps/webhook/start-webhook.sh

[Install]
WantedBy=multi-user.target